# roles/configs/files/autorestic/autorestic-prune.service
[Unit]
Description=Autorestic backup service (data pruning)
[Service]
Type=oneshot
ExecStart=autorestic forget --verbose --all --prune
